# scripts/macos_config/com.freelims.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.freelims.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/Shared/SDrive/freelims_production/start_production.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/Shared/SDrive/freelims_logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/SDrive/freelims_logs/launchd_stderr.log</string>
    <key>UserName</key>
    <string>shaun</string>
    <key>GroupName</key>
    <string>staff</string>
    <key>WorkingDirectory</key>
    <string>/Users/Shared/SDrive/freelims_production</string>
</dict>
</plist> 