# service_files/freelims-dev-backend.service
[Unit]
Description=FreeLIMS Development Backend
After=network.target postgresql.service

[Service]
Type=simple
User=shaun
WorkingDirectory=/Users/shaun/Documents/GitHub/projects/freelims
ExecStart=/Users/shaun/Documents/GitHub/projects/freelims/scripts/system/dev/run_dev_backend.sh
Restart=always
RestartSec=10
StandardOutput=append:/Users/shaun/Documents/GitHub/projects/freelims/logs/dev_backend.log
StandardError=append:/Users/shaun/Documents/GitHub/projects/freelims/logs/dev_backend_error.log

[Install]
WantedBy=multi-user.target
