# service_files/freelims-prod-frontend.service
[Unit]
Description=FreeLIMS Production Frontend
After=network.target freelims-prod-backend.service

[Service]
Type=simple
User=shaun
WorkingDirectory=/Users/shaun/Documents/GitHub/projects/freelims/frontend
ExecStart=/bin/bash -c 'cd /Users/shaun/Documents/GitHub/projects/freelims/frontend && cat > .env.production.local << ENVEOF
REACT_APP_API_URL=http://localhost:8002/api
PORT=3002
NODE_ENV=production
ENVEOF
&& npx serve -s build -l 3002'
Restart=always
RestartSec=10
StandardOutput=append:/Users/shaun/Documents/GitHub/projects/freelims/logs/prod_frontend.log
StandardError=append:/Users/shaun/Documents/GitHub/projects/freelims/logs/prod_frontend_error.log

[Install]
WantedBy=multi-user.target
